# config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Name=Polybar
Comment=Lightweight top bar
Exec=/home/stdk/.config/polybar/launch.sh
X-GNOME-Autostart-enabled=true
